# deploy/crowdquest-backup.service
[Unit]
Description=CrowdQuest PostgreSQL backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
Environment=CROWDQUEST_APP_DIR=/opt/crowdquest-current
Environment=CROWDQUEST_ENV_FILE=/etc/crowdquest/production.env
ExecStart=/usr/local/sbin/crowdquest-postgres-backup
User=root
Group=root
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/backups/crowdquest
